# src/rafthercal/rafthercal-web.service
[Unit]
Description=Rafthercal Web Config UI
After=network.target

[Service]
WorkingDirectory=/home/pi/rafthercal
ExecStart=/home/pi/rafthercal/venv/bin/python3 /home/pi/rafthercal/web_config.py
Environment=RAFTHERCAL_CONFIG_DIR=/home/pi/rafthercal
Restart=always
StandardOutput=journal
StandardError=journal
# Optional: set credentials via environment (remove if no auth wanted)
#Environment=WEB_USERNAME=admin
#Environment=WEB_PASSWORD=secret

[Install]
WantedBy=default.target
